# chat_ui/vibenanochat-chat.service
[Unit]
Description=VibeNanoChat Server
After=network.target

[Service]
Type=notify
User=colligo
Group=colligo
WorkingDirectory=/mnt/localssd/VibeNanoChat
Environment="PATH=/mnt/localssd/VibeNanoChat/.venv/bin:/usr/local/bin:/usr/bin"
Environment="LOG_LEVEL=info"
Environment="CHAT_SERVER_PORT=8003"

# Gunicorn with uvicorn workers for FastAPI/ASGI
ExecStart=/mnt/localssd/VibeNanoChat/.venv/bin/gunicorn \
    --config chat_ui/gunicorn_config.py \
    chat_ui.asgi:application

# Restart policy
Restart=always
RestartSec=10

# Security and resource limits
NoNewPrivileges=true
LimitNOFILE=65536
LimitNPROC=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vibenanochat

[Install]
WantedBy=multi-user.target
